# content/service/7/run
#!/bin/sh

cp -r /workdir/rclone.conf /mnt/data/config
